# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn app:app -w 8 -b localhost:7002 --timeout 1800
directory=/usr/src/app
