# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn app:app -w 10 -b localhost:5000
directory=/usr/src/app
